# configs/rc.local
# Bypass VPN for server IPs (must be BEFORE fwmark rule)
ip rule add to 46.8.233.202 lookup main priority 90
ip rule add to 31.77.77.47 lookup main priority 90
ip rule add to 94.126.153.10 lookup main priority 90

# VPN route after tun0 is ready
( for i in 1 2 3 4 5 6 7 8 9 10; do
    [ -e /sys/class/net/tun0 ] && break
    sleep 1
done; ip route replace table vpn default dev tun0 ) &

exit 0
